# private_dot_config/systemd/user/psd-resync.timer.d/override.conf
[Timer]
OnStartupSec=10min
OnUnitActiveSec=1h
RandomizedDelaySec=5min
