# private_dot_config/systemd/user/psd-resync.timer.d/override.conf
[Timer]
OnUnitActiveSec=30min
RandomizedDelaySec=5min
