# private_dot_config/systemd/user/i3-cycle-focus.service
[Unit]
Description=i3 cycle focus
PartOf=graphical-session.target

[Service]
ExecStart=/usr/local/bin/i3-cycle-focus --history 16

