# private_dot_config/systemd/user/mbsync.service
[Unit]
Description=Mailbox synchronization service

[Service]
Type=oneshot
Slice=nicy-cpu16.slice
ExecStart=/usr/bin/mbsync -c %E/mutt/mbsyncrc -qa
