# private_dot_config/systemd/user/picom.service
[Unit]
Description=Picom
Documentation=man:picom(1)
PartOf=graphical-session.target

[Service]
ExecStart=picom

