# private_dot_config/systemd/user/pkgs-check.service
[Unit]
Description=Check for package updates service

[Service]
Type=oneshot
Slice=nicy.slice
ExecStart=/usr/local/bin/pkgs check
