# private_dot_config/systemd/user/polybar@.service
[Unit]
Description=Polybar
Documentation=man:polybar(1) man:polybar(5)
PartOf=graphical-session.target

[Service]
ExecStart=polybar %i

