# private_dot_config/systemd/user/tvguide.service
[Unit]
Description=Update tvguide service

[Service]
Type=oneshot
Slice=cpu16.slice
ExecStart=/usr/bin/sh %h/bin/update-guidetv
