# private_dot_config/systemd/user/xsettingsd.service
[Unit]
Description=X settings daemon
Documentation=man:xsettingsd(1)
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/xsettingsd
